# ops/scripts/goldrush-backup.cron
# GoldRush — root cron entry for daily Postgres backup.
# Install with:
#     sudo cp /opt/goldrush/repo/ops/scripts/goldrush-backup.cron /etc/cron.d/goldrush-backup
#     sudo chmod 644 /etc/cron.d/goldrush-backup
# Cron picks it up automatically.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=""

# Daily at 03:00 UTC
0 3 * * * root /opt/goldrush/repo/ops/scripts/backup.sh >> /var/log/goldrush-backup.log 2>&1
